# supervisord.conf
[program:ssh]
priority=10
directory=/app/run.sh
command=/app/run.sh
autostart=true
autorestart=true
stopsignal=QUIT
